# private_dot_config/systemd/user/ludusavi-backup.timer
[Unit]
Description="Ludusavi backup timer"

[Timer]
OnCalendar=*-*-* 00:00:00
Unit=ludusavi-backup.service

[Install]
WantedBy=timers.target
