# private_dot_config/systemd/user/ludusavi-backup.service
[Unit]
Description="Ludusavi backup"

[Service]
ExecStart=flatpak run com.github.mtkennerly.ludusavi backup --force
